# ch340-multi-display.service
[Unit]
Description=CH340 USB listener (image + terminals on plug)
After=graphical-session.target

[Service]
Type=simple
ExecStart=/home/w/Documents/automatic_raspberry-pi/ch340-multi-display.sh
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
